# systemd/j1-noc-platform.service
[Unit]
Description=NexusCore NOC Platform (Docker Compose)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment="PROJECT_DIR=/opt/nexuscore"
Environment="LIVE_ENV=/etc/nexuscore/.env.live"

# Load live secrets into systemd environment
EnvironmentFile=/etc/nexuscore/.env.live

# On start: copy live env in, bring up compose, then scrub working .env
ExecStartPre=/bin/bash -c 'cp /etc/nexuscore/.env.live /opt/nexuscore/.env'
ExecStart=/usr/bin/docker compose up -d --build --remove-orphans
ExecStartPost=/bin/bash -c 'cp /opt/nexuscore/.env.example /opt/nexuscore/.env'

# On stop: stop containers, reset .env to safe template
ExecStop=/usr/bin/docker compose down
ExecStopPost=/bin/bash -c 'cp /opt/nexuscore/.env.example /opt/nexuscore/.env'

User=root
Group=root

[Install]
WantedBy=multi-user.target
